# qr3_17.conf
[program:qr3_17]
command=python3 /home/curator3/qr3_17/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_17/test.err.log
stdout_logfile=/home/curator3/qr3_17/test.err.log
